# modules/ambient_mozilla/templates/dev.consul.plist.erb
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.consul</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>run the consul agent all the time forever.</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/consul</string>
      <string>-data-dir=/opt/boxen/config/consul</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>UserName</key>
    <string><%= scope.lookupvar("::boxen_user") %></string>
  </dict>
</plist>
